# alignak-scheduler.ini
[program:alignak-scheduler]
command=/usr/bin/alignak-scheduler -c /etc/alignak/daemons/schedulerd.ini
autostart=true
user=nagios
log_stdout=false
log_stderr=true
logfile=/var/log/supervisor-alignak-scheduler.log
