# alignak-reactionner.ini
[program:alignak-broker]
command=/usr/bin/alignak-reactionner -c /etc/alignak/daemons/reactionnerd.ini
autostart=true
user=nagios
log_stdout=false
log_stderr=true
logfile=/var/log/supervisor-alignak-reactionner.log
