# alignak-broker.ini
[program:alignak-broker]
command=/usr/bin/alignak-broker -c /etc/alignak/daemons/brokerd.ini
autostart=true
user=nagios
log_stdout=false
log_stderr=true
logfile=/var/log/supervisor-alignak-broker.log
