# alignak-poller.ini
[program:alignak-poller]
command=/usr/bin/alignak-poller -c /etc/alignak/daemons/pollerd.ini
autostart=true
user=nagios
log_stdout=false
log_stderr=true
logfile=/var/log/supervisor-alignak-poller.log
