# alignak-arbiter.ini
[program:alignak-arbiter]
command=/usr/bin/alignak-arbiter -d -r -c /etc/alignak/daemons/arbiterd.ini -a /etc/alignak/alignak.cfg
autostart=true
user=nagios
log_stdout=false
log_stderr=true
logfile=/var/log/supervisor-alignak-arbiter.log
